# squirrel_interaction/autostart/image.desktop
[Desktop Entry]
Type=Application
Exec=eog --fullscreen /home/pi/catkin_ws/src/screen/media/images/look_front.png
